# server/weather/weather.cron
# m h  dom mon dow   command
40 * * * * rm -f /var/log/weatherjob.log; /opt/eink-display/weather/weather-script.sh >> /var/log/weatherjob.log 2>&1
10 * * * * rm -f /var/log/weatherjob.log; /opt/eink-display/weather/weather-script.sh >> /var/log/weatherjob.log 2>&1
